# supervisord.conf
; =============================================================================
; Musically — Supervisord Configuration
; =============================================================================
; Manages all services in a single self-contained container:
;   - PostgreSQL (database)
;   - Redis (message broker)
;   - Nginx (reverse proxy + static files)
;   - Uvicorn (FastAPI via start.sh)
;
; start.sh handles first-run DB init and Alembic migrations.
; All persistent data lives under /config (mounted volume).

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
childlogdir=/var/log/supervisor

; ---------------------------------------------------------------------------
; PostgreSQL — self-contained database
; ---------------------------------------------------------------------------
[program:postgres]
command=/bin/bash -c 'PG_BIN=$(ls -d /usr/lib/postgresql/*/bin 2>/dev/null | head -1); mkdir -p /config/postgres /run/postgresql && chown postgres:postgres /config/postgres /run/postgresql && exec su - postgres -c "$PG_BIN/postgres -D /config/postgres -c unix_socket_directories=/run/postgresql -c listen_addresses=127.0.0.1"'
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/postgres.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true
priority=1

; ---------------------------------------------------------------------------
; Redis — self-contained message broker
; ---------------------------------------------------------------------------
[program:redis]
command=/bin/bash -c "mkdir -p /config/redis && exec redis-server --dir /config/redis --save 900 1 --save 300 10 --save 60 10000 --appendonly yes --bind 127.0.0.1"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/redis.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true
priority=2

; ---------------------------------------------------------------------------
; Nginx — reverse proxy + React SPA static files
; ---------------------------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true
priority=10

; ---------------------------------------------------------------------------
; Uvicorn — FastAPI application server
; ---------------------------------------------------------------------------
; Calls start.sh which initializes the DB on first run, runs Alembic
; migrations, and launches uvicorn as the musically user.
[program:uvicorn]
command=/app/start.sh
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
redirect_stderr=true
priority=20
