# supervisord.conf
; =============================================================================
; Musically — Supervisord Configuration
; =============================================================================
; Manages nginx (reverse proxy + static files) and uvicorn (FastAPI) in the
; same container. The uvicorn program calls start.sh which runs Alembic
; migrations before starting the ASGI server.

[supervisord]
nodaemon=true
user=musically
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisor/supervisord.pid
childlogdir=/var/log/supervisor

; ---------------------------------------------------------------------------
; Nginx — reverse proxy + React SPA static files
; ---------------------------------------------------------------------------
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ---------------------------------------------------------------------------
; Uvicorn — FastAPI application server
; ---------------------------------------------------------------------------
; Calls start.sh which runs Alembic migrations before launching uvicorn.
[program:uvicorn]
command=/app/start.sh
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
